# manager/misc/supervisord.conf
[supervisord]
nodaemon=true
user=root

# log
[program:libvirtd]
command=/init-libvirtd.sh
startretries=10
redirect_stderr=true


[program:virtlockd]
command=/usr/sbin/virtlockd

[program:virtlogd]
command=/usr/sbin/virtlogd