# manager/misc/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0



[program:libvirtd]
command=/init-libvirtd.sh
autostart=true
autorestart=true
startretries=20
stopsignal=KILL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:virtlockd]
command=/usr/sbin/virtlockd

[program:virtlogd]
command=/usr/sbin/virtlogd