# manager/misc/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:libvirtd]
command=/init-libvirtd.sh
startretries=10

[program:virtlockd]
command=/usr/sbin/virtlockd

[program:virtlogd]
command=/usr/sbin/virtlogd