# utils/mows-vm-supervisor/image-builder/common/mows-agent-init.service
[Unit]
Description=MOWS agent init (per-boot setup driven by /mowsinit)
After=network-online.target sshd.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mows-agent-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
